# themes/luci-theme-tkth/root/etc/uci-defaults/luci-theme-tkth
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Think=/luci-static/tkth
	set luci.main.mediaurlbase=/luci-static/tkth
	commit luci
EOF
exit 0
